# webviewer.service
[Unit]
Description=WebViewer HTTPS 网站服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/webviewer
ExecStart=/root/.openclaw/workspace/webviewer/start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
